# infrastructure/docker/services/frontend/etc/service/apache2/run
#!/bin/sh

exec /usr/sbin/apache2ctl -DFOREGROUND